# luneos-session/lib/systemd/system/com.webos.service.peripheralmanager.service
[Unit]
Description=webos - "%n"
Requires=ls-hubd.service
After=ls-hubd.service

[Service]
Type=simple
OOMScoreAdjust=-500
ExecStart=/usr/sbin/peripheralmanager
Restart=on-failure
