# luneos-session/DEBIAN/postinst
#!/bin/bash
set -e

# Fix surface-manager.sh shebang (uses 'source' which is bash-specific, not POSIX sh)
STARTUP_SCRIPT="/usr/share/luna-surfacemanager/startup/surface-manager.sh"
if [ -f "$STARTUP_SCRIPT" ]; then
    if head -1 "$STARTUP_SCRIPT" | grep -q "^#!/bin/sh"; then
        sed -i '1s|#!/bin/sh|#!/bin/bash|' "$STARTUP_SCRIPT"
        echo "Fixed surface-manager.sh shebang for bash compatibility"
    fi
    # Ensure script is executable
    chmod +x "$STARTUP_SCRIPT"
fi

# Ensure all shell scripts in luna-surfacemanager startup are executable
if [ -d "/usr/share/luna-surfacemanager/startup" ]; then
    find /usr/share/luna-surfacemanager/startup -name "*.sh" -exec chmod +x {} \;
    echo "Made luna-surfacemanager startup scripts executable"
fi

# Ensure all shell scripts in LuneOS bin directories are executable
for dir in /usr/bin /usr/sbin /usr/libexec; do
    if [ -d "$dir" ]; then
        find "$dir" -maxdepth 1 -name "*.sh" -exec chmod +x {} \; 2>/dev/null || true
    fi
done

# Create ls-hubd PID directory
mkdir -p /var/run/ls2

# Create luneos group if it doesn't exist
if ! getent group luneos > /dev/null 2>&1; then
    groupadd -r luneos
    echo "Created luneos group"
fi

# Add all human users (UID >= 1000) to luneos group
for user in $(awk -F: '$3 >= 1000 && $3 < 65534 {print $1}' /etc/passwd); do
    if ! id -nG "$user" 2>/dev/null | grep -qw luneos; then
        usermod -aG luneos "$user" 2>/dev/null || true
        echo "Added $user to luneos group"
    fi
done

# Set LuneOS as the default display manager (takes effect on reboot)
# This is the proper Debian/Ubuntu way - no masking needed
echo "/usr/bin/luna-surfacemanager" > /etc/X11/default-display-manager

# Create display-manager.service symlink to luna-surfacemanager
# This tells systemd which DM to start at boot
ln -sf /usr/lib/systemd/system/luna-surfacemanager.service /etc/systemd/system/display-manager.service

# Enable Luna Service Hub (must start before surface-manager)
systemctl enable ls-hubd.service 2>/dev/null || true

# Enable db8 services
systemctl enable db8.service 2>/dev/null || true
systemctl enable db8-maindb.service 2>/dev/null || true

# Enable ConnMan (LuneOS network manager) - will take effect on reboot
systemctl enable connman.service 2>/dev/null || true

# Reload systemd to recognize new service links
systemctl daemon-reload

echo ""
echo "=============================================="
echo "LuneOS session installed successfully!"
echo ""
echo "IMPORTANT: Reboot to start LuneOS desktop."
echo "Your current session will continue normally."
echo "=============================================="

exit 0
