# luneos-session/DEBIAN/postinst
#!/bin/bash
set -e

# Fix surface-manager.sh shebang (uses 'source' which is bash-specific, not POSIX sh)
STARTUP_SCRIPT="/usr/share/luna-surfacemanager/startup/surface-manager.sh"
if [ -f "$STARTUP_SCRIPT" ]; then
    if head -1 "$STARTUP_SCRIPT" | grep -q "^#!/bin/sh"; then
        sed -i '1s|#!/bin/sh|#!/bin/bash|' "$STARTUP_SCRIPT"
        echo "Fixed surface-manager.sh shebang for bash compatibility"
    fi
fi

# Create ls-hubd PID directory
mkdir -p /var/run/ls2

# Create luneos group if it doesn't exist
if ! getent group luneos > /dev/null 2>&1; then
    groupadd -r luneos
    echo "Created luneos group"
fi

# Add all human users (UID >= 1000) to luneos group
for user in $(awk -F: '$3 >= 1000 && $3 < 65534 {print $1}' /etc/passwd); do
    if ! id -nG "$user" 2>/dev/null | grep -qw luneos; then
        usermod -aG luneos "$user" 2>/dev/null || true
        echo "Added $user to luneos group"
    fi
done

# Mask display managers (disable isn't enough - they can still be started by dependencies)
systemctl mask gdm3.service 2>/dev/null || true
systemctl mask gdm.service 2>/dev/null || true
systemctl mask lightdm.service 2>/dev/null || true
systemctl mask sddm.service 2>/dev/null || true

# Disable NetworkManager (LuneOS uses ConnMan)
systemctl disable NetworkManager.service 2>/dev/null || true
systemctl disable NetworkManager-wait-online.service 2>/dev/null || true

# Enable ConnMan
systemctl enable connman.service 2>/dev/null || true

# Enable Luna Service Hub (must start before surface-manager)
systemctl enable ls-hubd.service

# Enable db8 services
systemctl enable db8.service 2>/dev/null || true
systemctl enable db8-maindb.service 2>/dev/null || true

# Enable Luna Surface Manager
systemctl enable luna-surfacemanager.service

# Create symlink to graphical.target.wants
ln -sf /usr/lib/systemd/system/luna-surfacemanager.service /etc/systemd/system/graphical.target.wants/ 2>/dev/null || true

# Reload systemd
systemctl daemon-reload

echo "LuneOS session installed. Reboot to start LuneOS desktop."

exit 0
