# luneos-session/DEBIAN/postinst
#!/bin/bash
set -e

# Mask display managers (disable isn't enough - they can still be started by dependencies)
systemctl mask gdm3.service 2>/dev/null || true
systemctl mask gdm.service 2>/dev/null || true
systemctl mask lightdm.service 2>/dev/null || true
systemctl mask sddm.service 2>/dev/null || true

# Disable NetworkManager (LuneOS uses ConnMan)
systemctl disable NetworkManager.service 2>/dev/null || true
systemctl disable NetworkManager-wait-online.service 2>/dev/null || true

# Enable ConnMan
systemctl enable connman.service 2>/dev/null || true

# Enable Luna Service Hub (must start before surface-manager)
systemctl enable ls-hubd.service

# Enable db8 services
systemctl enable db8.service 2>/dev/null || true
systemctl enable db8-maindb.service 2>/dev/null || true

# Enable Luna Surface Manager
systemctl enable luna-surfacemanager.service

# Create symlink to graphical.target.wants
ln -sf /usr/lib/systemd/system/luna-surfacemanager.service /etc/systemd/system/graphical.target.wants/ 2>/dev/null || true

# Reload systemd
systemctl daemon-reload

echo "LuneOS session installed. Reboot to start LuneOS desktop."

exit 0
