# luneos-session/DEBIAN/prerm
#!/bin/bash
set -e

# Disable LuneOS services
systemctl disable luna-surfacemanager.service 2>/dev/null || true
systemctl disable ls-hubd.service 2>/dev/null || true
systemctl disable db8.service 2>/dev/null || true
systemctl disable db8-maindb.service 2>/dev/null || true
systemctl disable connman.service 2>/dev/null || true

# Restore GDM as the default display manager (if it exists)
if [ -f /usr/lib/systemd/system/gdm3.service ] || [ -f /usr/lib/systemd/system/gdm.service ]; then
    # Restore default-display-manager
    if [ -f /usr/lib/systemd/system/gdm3.service ]; then
        echo "/usr/sbin/gdm3" > /etc/X11/default-display-manager
        ln -sf /usr/lib/systemd/system/gdm3.service /etc/systemd/system/display-manager.service
    elif [ -f /usr/lib/systemd/system/gdm.service ]; then
        echo "/usr/sbin/gdm" > /etc/X11/default-display-manager
        ln -sf /usr/lib/systemd/system/gdm.service /etc/systemd/system/display-manager.service
    fi
    echo "Restored GDM as default display manager"
elif [ -f /usr/lib/systemd/system/lightdm.service ]; then
    echo "/usr/sbin/lightdm" > /etc/X11/default-display-manager
    ln -sf /usr/lib/systemd/system/lightdm.service /etc/systemd/system/display-manager.service
    echo "Restored LightDM as default display manager"
elif [ -f /usr/lib/systemd/system/sddm.service ]; then
    echo "/usr/bin/sddm" > /etc/X11/default-display-manager
    ln -sf /usr/lib/systemd/system/sddm.service /etc/systemd/system/display-manager.service
    echo "Restored SDDM as default display manager"
else
    # No known display manager found, just remove the symlink
    rm -f /etc/systemd/system/display-manager.service
    echo "No display manager found to restore"
fi

# Re-enable NetworkManager (Ubuntu's default)
systemctl enable NetworkManager.service 2>/dev/null || true

systemctl daemon-reload

echo ""
echo "LuneOS session removed. Reboot to return to your previous desktop."
echo ""

exit 0
